# deploy/production/systemd/teralya-backup-verify.service
[Unit]
Description=Verify latest Teralya PostgreSQL backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=deployer
Group=deployer
Environment=TERALYA_ROOT=/srv/teralya
ExecStart=/srv/teralya/scripts/verify-latest-backup.sh

[Install]
WantedBy=multi-user.target
