# deploy/production/systemd/teralya-backup-verify.timer
[Unit]
Description=Weekly Teralya PostgreSQL backup verification

[Timer]
OnCalendar=Sun *-*-* 04:15:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=teralya-backup-verify.service

[Install]
WantedBy=timers.target
