# deploy/production/systemd/teralya-backup.timer
[Unit]
Description=Daily Teralya PostgreSQL backup

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=teralya-backup.service

[Install]
WantedBy=timers.target
